# themes/common/theme/dotfiles/.config/systemd/user/check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/alba/.config/scripts/battery/check-battery.sh
